# ferex_serial.conf
[program:ferex_serial]
command=python ferex_serial.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ferex_serial.log
