# gps_serial.conf
[program:gps_serial]
command=python gps_serial.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gps_serial.log
